# yausbir_lirc.service
[Unit]
Description=Small UDP LIRC daemon for yausbir

[Service]
Type=forking
ExecStart=/usr/sbin/yausbir_lirc

[Install]
WantedBy=multi-user.target
